# cron-backup/config/crontab
# 毎日3時にバックアップを行う
0 3 * * * flock -n /tmp/echo.lock /var/backup_sh/pg_dumpall.sh
# 30分に1回ごとバックアップ
*/30 * * * * flock -n /tmp/echo.lock /var/backup_sh/pg_rman.sh

